# local.d/backlight_brightness.start
#!/bin/sh
echo 290 > /sys/class/backlight/intel_backlight/brightness
